# scripts/qos_daemon.service
[Unit]
Description=QOS Daemon - JW Attendant Scheduler Monitoring Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/jw-attendant-scheduler-django
ExecStart=/usr/bin/python3 /opt/jw-attendant-scheduler-django/scripts/qos_daemon.py /opt/jw-attendant-scheduler-django /opt/jw-attendant-scheduler-django/scripts/qos_config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qos-daemon

# Environment variables
Environment=PYTHONPATH=/opt/jw-attendant-scheduler-django
Environment=DJANGO_SETTINGS_MODULE=jw_scheduler.settings

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/jw-attendant-scheduler-django/logs
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
